# systemd/timelapse-render.service
[Unit]
Description=Timelapse Render Worker
After=network.target
PartOf=timelapse.target

[Service]
Type=simple
User=pls
Group=pls
ExecStart=/home/pls/timelapse-project/timelapse/.venv/bin/timelapse run render --config /etc/timelapse/timelapse.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/usb/timelapse

[Install]
WantedBy=timelapse.target
